# fs/etc/systemd/system/amazon-cloudwatch-agent.service.d/00_fetch-config.conf
[Unit]
Wants=aws-env.service
After=aws-env.service

[Service]
ExecStartPre=/bin/sh -c 'systemctl is-active aws-env && /usr/libexec/amazon-cloudwatch-agent/fetch-config.sh'
